# profile.d/google.sh
export GOMA_DIR=$HOME/goma
export CHROME_SRC=$HOME/src/chromium/src

export PATH=$PATH:$HOME/src/depot_tools
