# profile.d/google.sh
export GOMA_DIR=$HOME/goma
export CHROME_SRC=$HOME/src/chromium/src
export CHROMEOS_SRC=$HOME/src/chromiumos

export PATH=$PATH:$HOME/src/depot_tools
